# installers/macos/pkg-scripts/postinstall
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
#
# postinstall — Run after .pkg installs files to ~/.coding-agents-kit/
#
# Arguments from macOS Installer:
#   $1 = package path
#   $2 = install location (e.g., /Users/username for user-domain installs)
#   $3 = target volume
#   $4 = root directory for the system

if [[ -z "$2" ]] || [[ ! -d "$2" ]]; then
    echo "postinstall: invalid install location: $2" >&2
    exit 1
fi

# Installer passes "$2" with a trailing /. for user-domain installs; strip
# textually rather than via cd/pwd -P so we keep the conventional path
# (preserve a symlinked HOME instead of resolving it).
INSTALL_LOCATION="${2%/.}"
INSTALL_LOCATION="${INSTALL_LOCATION%/}"
PREFIX="${INSTALL_LOCATION}/.coding-agents-kit"
USER_HOME="${INSTALL_LOCATION}"

# ---------------------------------------------------------------------------
# Render launchd templates
# ---------------------------------------------------------------------------

PLIST_DIR="${USER_HOME}/Library/LaunchAgents"
PLIST_FILE="${PLIST_DIR}/dev.falcosecurity.coding-agents-kit.plist"

mkdir -p "$PLIST_DIR"
mkdir -p "$PREFIX/log"

# Render plist with actual prefix and HOME.
sed -e "s|@PREFIX@|${PREFIX}|g" -e "s|@HOME@|${USER_HOME}|g" \
    "$PREFIX/launchd/dev.falcosecurity.coding-agents-kit.plist" \
    > "$PLIST_FILE"

# Render launcher script.
sed "s|@PREFIX@|${PREFIX}|g" \
    "$PREFIX/launchd/coding-agents-kit-launcher.sh" \
    > "$PREFIX/bin/coding-agents-kit-launcher.sh"
chmod 755 "$PREFIX/bin/coding-agents-kit-launcher.sh"

# ---------------------------------------------------------------------------
# Create user rules directory (if not exists)
# ---------------------------------------------------------------------------

if [[ ! -d "$PREFIX/rules/user" ]]; then
    mkdir -p "$PREFIX/rules/user"
fi

# ---------------------------------------------------------------------------
# Start service
# ---------------------------------------------------------------------------
#
# Postinstall runs either:
#   - as root, when the user invokes `sudo installer -pkg …` or authenticates
#     in the GUI Installer for a system-scoped component; OR
#   - as the installing user, when the pkg is user-domain (our case, via
#     `enable_currentUserHome="true"`) and the GUI Installer.app or
#     `installer -target CurrentUserHomeDirectory` runs the scripts without
#     elevation.
#
# In the elevated case we can't just call `launchctl load` from root — the
# agent would land in the `system`/`background` domain, not the user's
# `gui/$UID` session where the GUI clients expect it. `launchctl asuser $UID`
# explicitly crosses into the user's launchd session and is the
# Apple-documented way to load a LaunchAgent from a root context.
#
# Errors are NOT redirected to /dev/null so Installer.app surfaces them in
# "Show Details"; `|| true` keeps the overall install from failing if the
# agent can't load for an unrelated reason (the user can retry with ctl).

OWNER=$(stat -f '%Su' "$USER_HOME")
OWNER_UID=$(id -u "$OWNER" 2>/dev/null || echo "")

if [[ "$(id -u)" -eq 0 && -n "$OWNER_UID" && "$OWNER" != "root" ]]; then
    # When installer runs as root, PackageKit extracts the payload with root
    # ownership even for user-domain targets. Reassign the install tree and
    # rendered plist to the installing user so they can manage the service
    # and uninstall without sudo.
    chown -R "$OWNER" "$PREFIX" || true
    chown "$OWNER" "$PLIST_FILE" || true
    launchctl asuser "$OWNER_UID" /bin/launchctl load "$PLIST_FILE" || true
else
    launchctl load "$PLIST_FILE" || true
fi

exit 0
